# .config/systemd/user/mbsync@.timer
[Unit]
Description=Runs mbsync each timeout

[Timer]
OnUnitInactiveSec=30

[Install]
WantedBy=xinitrc.target
