# .config/systemd/user/mbsync@.service
[Unit]
Description=Start mbsync for mbox: %i

[Service]
Type=simple
ExecStart=/usr/bin/mbsync -q %i

[Install]
WantedBy=xinitrc.target
