# .config/systemd/user/xinitrc.target
[Unit]
Description=Xinitrc Stuff
After=environment.target
Wants=environment.target

[Install]
Alias=default.target
